# init/pqpmd.service
[Unit]
Description=PQPM Process Manager Daemon
Documentation=https://github.com/pqpm/pqpm
After=network.target local-fs.target

[Service]
Type=simple
# Give the system a bit more time to mount home directories if needed
ExecStartPre=/usr/bin/sleep 2
ExecStart=/usr/local/bin/pqpmd
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pqpmd

# Runtime directories (removed automatic systemd management to fix namespace error)
# The daemon will handle directory creation on startup.

# Security hardening (all isolation disabled for troubleshooting)
# Using 'false' instead of 'off' or 'no' for consistency
ProtectSystem=false
ProtectHome=false
PrivateTmp=false
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
